# deploy/h1scope.service
[Unit]
Description=H1 Scope Dashboard Backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/h1scope/backend
ExecStart=/var/www/h1scope/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8765 --workers 1
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/var/www/h1scope/backend/.env

[Install]
WantedBy=multi-user.target
